# core/php80/crontab/jobs
* * * * * /usr/local/bin/php /var/www/html/php80/artisan schedule:run >> /dev/null 2>&1
