# core/php74/crontab/jobs
* * * * * /usr/local/bin/php /var/www/html/php74/artisan schedule:run >> /dev/null 2>&1
